# files/common/etc/profile.d/exorint.sh
#!/bin/sh

if [ "$USER" = "admin" ]; then
    PATH=$PATH:/usr/local/sbin:/usr/sbin:/sbin
fi

umask 002
