# files/common/etc/network/if-up.d/00-fix-metrics
. /etc/exorint.funcs

exorint_iface_metric_apply $IFACE
