# files/common/etc/network/if-up.d/10-dns-restart
#!/bin/sh

/etc/init.d/dnsmasq restart

exit 0
